# systemd/meteor-scatter@.service
[Unit]
Description=meteor-scatter MSK144 meteor-scatter spot recorder for radiod %I
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=10

[Service]
Type=notify
User=meteorscat
Group=meteorscat

WatchdogSec=120
NotifyAccess=main

EnvironmentFile=-/etc/sigmond/coordination.env
# %i (literal instance) is the working key: %I "unescapes" systemd
# escaping, turning dashes into slashes (B4-100-rx888mk2 →
# B4/100/rx888mk2), so the %I path never matches a real flat file.
# Same bug + fix as wspr-recorder@.service.  %I kept second for any
# host that worked around it with nested directories.
EnvironmentFile=-/etc/meteor-scatter/env/%i.env
EnvironmentFile=-/etc/meteor-scatter/env/%I.env

WorkingDirectory=/opt/git/sigmond/meteor-scatter

ExecStartPre=+/bin/mkdir -p /var/lib/meteor-scatter/%I/msk144
ExecStartPre=+/bin/mkdir -p /var/log/meteor-scatter
ExecStartPre=+/usr/bin/chown -R meteorscat:meteorscat /var/lib/meteor-scatter/%I
ExecStartPre=+/usr/bin/chown -R meteorscat:meteorscat /var/log/meteor-scatter
# /var/lib/hs-uploader is SHARED across HamSCI client users (meteorscat,
# wsprdaemon, ...) — its setgid+group-write shape is owned by
# tmpfiles.d/hs-uploader.conf (root:sigmond, mode 02775).  Earlier
# versions of this unit had:
#
#   ExecStartPre=+/bin/mkdir -p /var/lib/hs-uploader
#   ExecStartPre=+/usr/bin/chown meteorscat:meteorscat /var/lib/hs-uploader
#
# Both wrong: the first redundant with tmpfiles.d, the second hostile —
# it re-stomped the dir's group from sigmond to meteorscat on every
# meteor-scatter restart, locking wsprdaemon-client's wd-upload-hs out
# of the shared watermarks.db with "attempt to write a readonly
# database".  Observed on bee1 2026-05-14.  Removed.

# Phase-3 cutover (sigmond MULTI-INSTANCE-ARCHITECTURE.md):
# --instance %I picks up /etc/meteor-scatter/%I.toml when the operator
# has run `sudo smd instance migrate` (Phase 8).  Until then, the
# per-instance config doesn't exist and the daemon falls back to the
# legacy /etc/meteor-scatter/meteor-scatter-config.toml with a one-line
# deprecation warning.  --radiod-id %I stays for the legacy fallback
# path (selects the matching [[radiod]] block within the shared
# config); ignored once per-instance config takes over.
# Note: --config is intentionally omitted — passing it would force
# the legacy path always (explicit_path wins in resolve_config_path).
ExecStart=/opt/git/sigmond/meteor-scatter/venv/bin/python3 -m meteor_scatter.cli daemon \
    --instance %I \
    --radiod-id %I

# 180s is generous for the 2-channel (10 m + 6 m) MSK144 layout: each
# channel triggers a setpolicy round-trip with radiod, the initial
# multicast discovery sleeps 1s, and a fresh start after a stop adds
# ~10s of teardown.  Kept at 180s (same as the wspr/psk siblings) for
# headroom on a busy host even though only two channels provision here.
TimeoutStartSec=180

Restart=always
# EX_CONFIG (78): daemon hit the unconfigured radiod placeholder and exited
# fast (see cli._handle_daemon guard). Do NOT crash-loop a config that can
# never succeed — stop cleanly until an operator runs `config init`.
RestartPreventExitStatus=78
RestartSec=5

# Logs go to the systemd journal (sigmond standard) — query with
# `journalctl -u meteor-scatter@%I` or `smd log meteor-scatter`.
StandardOutput=journal
StandardError=inherit
SyslogIdentifier=meteor-scatter@%I

Nice=5
# jt9 is a Fortran binary that mmaps its own working buffers per slot
# (~60 MB RSS each).  MSK144 provisions only two channels (10 m + 6 m),
# so at most two concurrent jt9 children — far under this ceiling.  1 GB
# leaves comfortable headroom for those plus the meteor-scatter Python
# process itself; kept generous so a slow decode never trips the cgroup.
MemoryMax=1G
MemorySwapMax=0

NoNewPrivileges=true
PrivateTmp=false
ProtectSystem=strict
ProtectHome=read-only
# /var/lib/sigmond holds the shared hamsci_sink queue (sink.db); without
# it RW, ProtectSystem=strict makes the dir read-only and the psk Writer
# silently no-ops (every decoded spot dropped before it reaches the queue).
ReadWritePaths=/var/lib/meteor-scatter /var/log/meteor-scatter /var/lib/hs-uploader /var/lib/sigmond
ReadOnlyPaths=/etc/meteor-scatter
PrivateNetwork=false

CapabilityBoundingSet=CAP_NET_RAW CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_NET_RAW CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
